# launchd/com.nkmr.issues-autobackup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nkmr.issues-autobackup</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nkmr/bin/git-auto-backup.sh</string>
    <string>/Users/nkmr/ghq/github.com/nkmr-jp/issues</string>
    <string>--llm</string>
  </array>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/nkmr/Library/Logs/issues-autobackup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nkmr/Library/Logs/issues-autobackup.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/nkmr</string>
  </dict>
</dict>
</plist>
