# launchd/com.nkmr.claude-stall-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nkmr.claude-stall-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/nkmr/bin/claude-stall-monitor.sh</string>
    <string>--watch</string>
  </array>
  <key>StartInterval</key>
  <integer>30</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/nkmr/Library/Logs/claude-stall-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/nkmr/Library/Logs/claude-stall-monitor.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/nkmr</string>
  </dict>
</dict>
</plist>
